# scripts/controlblock.service
[Unit]
Description=ControlBlock Service

[Service]
Type=oneshot
ExecStart=/usr/bin/controlblock

[Install]
WantedBy=multi-user.target